# crontab
* * * * * aws s3 sync s3://demo15544 /var/www/html
